# 3/supervisord_shadowsocks.conf
[program:socks]
command=ssserver -c /etc/shadowsocks.json
startsecs=0
autorestart=true
