# 2/zzsupervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:lantern]
command=/usr/bin/lantern --addr 0.0.0.0:8787 

[program:ssserver]
command=/usr/bin/proxychains4 ssserver -p 443 -k wuli.ta0ta0 -d start
